# core/templates/api-server-launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hermes.api.{{PROFILE}}</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{HERMES_BIN}}</string>
        <string>--profile</string>
        <string>{{PROFILE}}</string>
        <string>gateway</string>
        <string>run</string>
        <string>--replace</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{HERMES_HOME}}</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>{{HOME}}</string>
        <key>HERMES_HOME</key>
        <string>{{HERMES_HOME}}</string>
        <key>HERMES_PROFILE</key>
        <string>{{PROFILE}}</string>
        <key>API_SERVER_PORT</key>
        <string>{{PORT}}</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/tmp/api-{{PROFILE}}-{{PORT}}.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/api-{{PROFILE}}-{{PORT}}.err</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
