# s6m-config/plists/com.hermes.a2a.registry.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hermes.a2a.registry</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/alexcai/.hermes/plugins/hermes-a2a/gateway.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/alexcai/.hermes/plugins/hermes-a2a</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>GATEWAY_HOST</key>
        <string>127.0.0.1</string>
        <key>GATEWAY_PORT</key>
        <string>8928</string>
        <key>PORT_MAP_PATH</key>
        <string>/Users/alexcai/code/hermes-a2a/s6m-config/port-map.md</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/alexcai/.hermes/logs/a2a-registry.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/alexcai/.hermes/logs/a2a-registry.err</string>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
